# root/etc/services.d/update-feeds/run
#!/usr/bin/with-contenv bash

[[ -d /var/www/html ]] && \
        s6-setuidgid abc php \
        /var/www/html/update_daemon2.php > /dev/null 2>&1
