# root/etc/services.d/update-feeds/run
#!/usr/bin/with-contenv bash

[[ -d /app/www/public ]] && \
        s6-setuidgid abc php \
        /app/www/public/update_daemon2.php > /dev/null 2>&1
